# nginx-fpm/php7.2.12/supervisord.conf
[supervisord]
nodaemon=true
logfile=/home/LogFiles/supervisor/supervisord.log

[program:sshd]
command=rc-service sshd start

[program:nginx]
command=sh /usr/local/bin/super_nginx.sh